# server/nginx/record.conf
description "Gunicorn application server running record"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
setuid ubuntu
setgid www-data

chdir /home/ubuntu/Record-Collection/server
exec venv/bin/gunicorn --workers 3 --bind 127.0.0.1:8001 collection.wsgi:application
